# config/systemd/user/waybar-weather.timer
[Unit]
Description=Refresh weather once daily

[Timer]
OnBootSec=30s
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
