# home/.config/systemd/user/waybar-weather.service
[Unit]
Description=Fetch weather for Waybar

[Service]
Type=oneshot
ExecStart=%h/.config/waybar/scripts/weather-fetch
